# wnsm-smartmeter/rootfs/etc/services.d/wnsm/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle service shutdown
# ==============================================================================
bashio::log.info "Wiener Netze Smart Meter service terminated"